# 集群配置&脚本/伪集群版/jenkins/jenkinsService.sh
#! /bin/bash

# chkconfig:   2345 90 10

# description:  jenkinsService



java -jar /opt/jenkins.war



