# 集群配置&脚本/分布式集群版/kafka/kafka-eagle/kafka-eagle.service
[Unit]
Description=kafka监控服务
After=network.target
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/local/kafka/script/startkafkaEagle.sh start

[Install]
WantedBy=multi-user.target
